# contrib/legacy-ubuntu-init/etc/init/zun-compute.conf
description "OpenStack Container Service Compute Agent"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --chuid zun \
--chdir /var/lib/zun \
--name zun-compute \
--exec /usr/local/bin/zun-compute -- \
--config-file=/etc/zun/zun.conf
