# contrib/legacy-ubuntu-init/etc/init/zun-api.conf
description "OpenStack Container Service API"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec start-stop-daemon --start --chuid zun \
--chdir /var/lib/zun \
--name zun-api \
--exec /usr/local/bin/zun-api -- \
--config-file=/etc/zun/zun.conf
